# gear_server.sh
#!/bin/sh

# PROVIDE: gear_server

. /etc/rc.subr

name="gear_server"
